# ops/runtime-ownership/blackspire-command.service
# Blackspire Command — durable production runtime (REVIEW TEMPLATE, NOT INSTALLED)
#
# Reviewed systemd unit template. NOT installed or enabled by any repository change. Applying
# it is an operator action under blocker #2 approval. Runs the production supervisor
# (which calls verifyVpsRuntime and fails closed before spawning api/worker) as the
# non-root blackspire user from the immutable current release.

[Unit]
Description=Blackspire Command control plane (durable production runtime)
After=network-online.target
Wants=network-online.target
# Restart-storm cap: modern systemd honors StartLimit* only in [Unit] (not [Service]). At
# five failed starts within ten minutes the unit enters failed state and stops for operator
# review instead of looping. Paired with Restart=on-failure/RestartSec=5 in [Service].
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple

# --- Least-privileged identity (blocker #2) ---
User=blackspire
Group=blackspire

# --- Release + persistent state layout ---
WorkingDirectory=/opt/blackspire-command/current
# Secrets live outside Git in a 0640 root:blackspire file; values never appear in the repo.
EnvironmentFile=/etc/blackspire/command.env

# Pinned PATH. systemd's manager PATH is /usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/snap/bin
# and excludes /opt/nodejs, so any helper resolving `node` through PATH would silently get the
# distribution's Node 18, which lacks node:sqlite. The reviewed interpreter directory is placed
# first so no PATH-resolved lookup inside ExecStartPre can reach /usr/bin/node. The helpers also
# resolve the interpreter themselves via scripts/lib/node-bin.sh, so this is defense in depth
# rather than the only guard.
Environment=PATH=/opt/nodejs/node-v22.23.1-linux-x64/bin:/usr/local/bin:/usr/bin:/bin

# Verify the environment profile before starting; this also enforces the loopback-only
# BIND_HOST, the explicit non-conflicting PORT, and read-only port-conflict detection. The
# supervisor re-verifies both at runtime and fails closed on an occupied port.
ExecStartPre=/usr/bin/env bash scripts/verify-environment.sh vps-production
# The interpreter is pinned to an absolute path on purpose. systemd's manager PATH is
# /usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/snap/bin and does NOT include /opt/nodejs,
# so `/usr/bin/env node` resolves to the distribution's /usr/bin/node (18.x on this host). The
# control plane requires node:sqlite (packages/task-engine/db.js, scripts/backup.js,
# scripts/restore.js), which does not exist before Node 22.5, so an unpinned unit fails every
# start, exhausts StartLimitBurst, and leaves the unit failed. The restricted staging unit
# already pins this same absolute path; production must match it.
ExecStart=/opt/nodejs/node-v22.23.1-linux-x64/bin/node scripts/production-supervisor.js

# --- Restart policy: cap at five attempts in ten minutes, then stop for operator review ---
# (StartLimitIntervalSec/StartLimitBurst live in [Unit] above — that is where systemd reads them.)
Restart=on-failure
RestartSec=5

# --- Sandboxing / hardening ---
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
# Only the persistent state tree is writable; releases and everything else stay read-only.
ReadWritePaths=/opt/blackspire-command/shared
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
RestrictNamespaces=yes
# The production port is unprivileged (>= 1024) and loopback-only, so no capabilities are
# needed to bind it. BIND_HOST=127.0.0.1 and the explicit PORT come from the EnvironmentFile;
# ports 8787 (existing API/worker) and 8788 (restricted staging) are rejected by the preflight.
CapabilityBoundingSet=
AmbientCapabilities=
# NOTE: MemoryDenyWriteExecute is intentionally NOT set — the Node/V8 JIT requires W^X-exempt
# executable memory and would crash under it.

[Install]
WantedBy=multi-user.target
