# ops/runtime-ownership/blackspire-command.service
# Blackspire Command — durable production runtime (REVIEW TEMPLATE, NOT INSTALLED)
#
# Reviewed systemd unit template. NOT installed or enabled by any repository change. Applying
# it is an operator action under blocker #2 approval. Runs the production supervisor
# (which calls verifyVpsRuntime and fails closed before spawning api/worker) as the
# non-root blackspire user from the immutable current release.

[Unit]
Description=Blackspire Command control plane (durable production runtime)
After=network-online.target
Wants=network-online.target
# Restart-storm cap: modern systemd honors StartLimit* only in [Unit] (not [Service]). At
# five failed starts within ten minutes the unit enters failed state and stops for operator
# review instead of looping. Paired with Restart=on-failure/RestartSec=5 in [Service].
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple

# --- Least-privileged identity (blocker #2) ---
User=blackspire
Group=blackspire

# --- Release + persistent state layout ---
WorkingDirectory=/opt/blackspire-command/current
# Secrets live outside Git in a 0640 root:blackspire file; values never appear in the repo.
EnvironmentFile=/etc/blackspire/command.env

# Verify the environment profile before starting; this also enforces the loopback-only
# BIND_HOST, the explicit non-conflicting PORT, and read-only port-conflict detection. The
# supervisor re-verifies both at runtime and fails closed on an occupied port.
ExecStartPre=/usr/bin/env bash scripts/verify-environment.sh vps-production
ExecStart=/usr/bin/env node scripts/production-supervisor.js

# --- Restart policy: cap at five attempts in ten minutes, then stop for operator review ---
# (StartLimitIntervalSec/StartLimitBurst live in [Unit] above — that is where systemd reads them.)
Restart=on-failure
RestartSec=5

# --- Sandboxing / hardening ---
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
# Only the persistent state tree is writable; releases and everything else stay read-only.
ReadWritePaths=/opt/blackspire-command/shared
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
RestrictNamespaces=yes
# The production port is unprivileged (>= 1024) and loopback-only, so no capabilities are
# needed to bind it. BIND_HOST=127.0.0.1 and the explicit PORT come from the EnvironmentFile;
# ports 8787 (existing API/worker) and 8788 (restricted staging) are rejected by the preflight.
CapabilityBoundingSet=
AmbientCapabilities=
# NOTE: MemoryDenyWriteExecute is intentionally NOT set — the Node/V8 JIT requires W^X-exempt
# executable memory and would crash under it.

[Install]
WantedBy=multi-user.target
